# .config/update-motd.d/51-textbanner
#!/bin/bash

lsb_release -ds | figlet -t -f digital

