# .config/update-motd.d/59-weather
#!/bin/bash

ansiweather -s true -l ambares-et-lagrave

